# archlinux/etc/systemd/system/auditor.timer
[Unit]
Description=Run auditor.service every week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
