# archlinux/etc/systemd/system/auditor.service
[Unit]
Description=Run the auditor script to check for system misconfigurations

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auditor
