# archlinux/etc/systemd/system/pacman-notify.timer
[Unit]
Description=Run pacman-notify.service every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
