# rootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
Type=simple
Restart=on-failure
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --skip-login --nonewline --noissue --noclear --autologin username_placeholder %I $TERM
