# rootfs/etc/systemd/user/journalctl-notify.service
[Unit]
Description=Follow journalctl and send desktop notification accordingly
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/journalctl-notify

[Install]
WantedBy=graphical-session.target

