# archlinux/etc/systemd/user/restic-unattended.timer
[Unit]
Description=Run restic-unattended.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

