# rootfs/etc/systemd/system/btrfs-balance.timer
[Unit]
Description=Run btrfs-balance.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
