# rootfs/etc/systemd/system/pacman-notify.timer
[Unit]
Description=Run pacman-notify.service every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
