# rootfs/etc/systemd/system/pacman-sync.timer
[Unit]
Description=Run pacman-sync.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
