# rootfs/etc/systemd/system/pacman-notify.service
[Unit]
Description=Send desktop notifications related to pacman

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pacman-notify
