# rootfs/etc/systemd/system/should-reboot-check.timer
[Unit]
Description=Run should-reboot-check.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
