# rootfs/etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
# Add --any so that the network is considered online if any of the interface is ready
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any
